# files/etc/init.d/wstunnel
#!/bin/sh /etc/rc.common
# OpenWrt init.d script for wstunnel server

START=99
STOP=99
USE_PROCD=1

# 定义二进制文件路径
WSTUNNEL_BIN="/usr/bin/wstunnel"

start_service() {
    # 检查二进制文件是否存在
    [ -x "$WSTUNNEL_BIN" ] || {
        echo "Error: $WSTUNNEL_BIN not found or not executable"
        return 1
    }

    # 读取 UCI 配置
    config_load wstunnel
    config_get enabled server enabled '0'
    config_get listen_addr server listen_addr 'wss://0.0.0.0:443'
    config_get secret server secret 'secret'
    config_get restrict_to server restrict_to 'localhost:51820'
    config_get log_level server log_level 'INFO'
    config_get tls_cert server tls_certificate ''
    config_get tls_key server tls_private_key ''

    # 检查是否启用
    [ "$enabled" -eq 0 ] && return 1

    # 构建启动参数
    local args="server"
    args="$args --restrict-http-upgrade-path-prefix $secret"
    args="$args --restrict-to $restrict_to"
    args="$args --log-lvl $log_level"

    # 如果提供了自定义 TLS 证书和密钥，添加到参数中
    [ -n "$tls_cert" ] && [ -n "$tls_key" ] && {
        args="$args --tls-certificate $tls_cert"
        args="$args --tls-private-key $tls_key"
    }

    # 最后添加监听地址
    args="$args $listen_addr"

    # 使用 procd 启动守护进程
    procd_open_instance
    procd_set_param command "$WSTUNNEL_BIN"
    procd_append_param command $args
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-5}
    procd_close_instance

    echo "wstunnel server started with args: $args"
}

stop_service() {
    # procd 会自动处理停止
    echo "Stopping wstunnel server..."
}